# deploy/proxmox/gutafinn-backup.service
[Unit]
Description=Gutafinn SQLite backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/gutafinn/deploy/proxmox
ExecStart=/opt/gutafinn/deploy/proxmox/backup.sh
