# deploy/proxmox/gutafinn-backup.timer
[Unit]
Description=Nightly backup for Gutafinn

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=gutafinn-backup.service

[Install]
WantedBy=timers.target
